# support-files/mariadb-extra@.socket.in

[Unit]

Description=MariaDB @VERSION@ database server (socket activation extra port multi-instance %I)
Documentation=man:mariadbd(8)
Documentation=https://mariadb.com/docs/server/server-management/starting-and-stopping-mariadb/systemd

[Install]
WantedBy=sockets.target

[Socket]

Service=mariadb@%i.service

# An "extra" as a descriptor name that means treat these ListenStreams as
# the same as an extra_port.
# Ref: https://mariadb.com/docs/server/ha-and-performance/optimization-and-tuning/buffers-caches-and-threads/thread-pool/thread-pool-system-status-variables

FileDescriptorName=extra

ListenStream=@mariadb-extra-%I
ListenStream=@MYSQL_UNIX_ADDR@-extra-%I
